# ansible/roles/llm-server/templates/llama-server.service.j2
[Unit]
Description=llama-server {{ llm_server_model_alias | default('ik_llama.cpp') }}
After=network.target
# Uncomment if nvidia-persistenced is installed
# After=network.target nvidia-persistenced.service

{% set selected_profile = llm_server_profiles[llm_server_profile] %}
[Service]
Type=simple
User={{ llm_server_user }}
WorkingDirectory={{ llm_server_install_dir }}/ik_llama.cpp
EnvironmentFile={{ llm_server_environment_file }}
ExecStartPre=/usr/bin/test -f ${LLAMA_MODEL}
ExecStart=/usr/bin/numactl --interleave=all \
    {{ llm_server_install_dir }}/ik_llama.cpp/build/bin/llama-server \
    --model ${LLAMA_MODEL} \
    --alias ${LLAMA_ALIAS} \
    -ngl ${LLAMA_GPU_LAYERS} \
{% for ot in llm_server_tensor_overrides | default([]) %}
    -ot "{{ ot.pattern }}={{ ot.backend }}" \
{% endfor %}
    --split-mode ${LLAMA_SPLIT_MODE} \
    --tensor-split {{ llm_server_tensor_split }} \
    --jinja \
    --ctx-size ${LLAMA_CTX_SIZE} \
    --cache-type-k {{ llm_server_cache_type_k }} \
    --cache-type-v {{ llm_server_cache_type_v }} \
{% if llm_server_threads | int > 0 %}
    --threads {{ llm_server_threads }} \
{% endif %}
{% if llm_server_threads_batch | int > 0 %}
    --threads-batch {{ llm_server_threads_batch }} \
{% endif %}
    --host {{ llm_server_host }} \
    --port {{ llm_server_port }} \
    --parallel {{ llm_server_parallel }} \
    --reasoning-format {{ selected_profile.reasoning_format }} \
    --temp {{ selected_profile.temperature }} \
    --top-p {{ selected_profile.top_p }} \
    --top-k {{ selected_profile.top_k }}
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=3
LimitNOFILE=65536
# Optional: lock model and KV cache in memory to reduce latency jitter
# Enable by uncommenting and adding --mlock --no-mmap to ExecStart
# LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
